# launchd/com.ralph.daily-compound-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ralph.daily-compound-review</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PROJECT_DIR__/scripts/daily-compound-review.sh</string>
    <string>__PROJECT_DIR__</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__PROJECT_DIR__</string>

  <!-- Run at 10:30 PM every night -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>22</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>__PROJECT_DIR__/logs/compound-review.log</string>

  <key>StandardErrorPath</key>
  <string>__PROJECT_DIR__/logs/compound-review.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>__HOME_DIR__</string>
  </dict>
</dict>
</plist>
